# sdcard-backup.path
[Unit]
Description=Watch for camera SD Card Mount

[Path]
PathChanged=/media/kyla/3232-6238/DCIM
Unit=sdcard-backup.service

[Install]
WantedBy=multi-user.target
